# cozy-controller.conf
[program:cozy-controller]
autorestart=true
command=cozy-controller
environment=NODE_ENV="production"
redirect_stderr=true
user=root
