# supervisor/cozy-init.conf
[program:cozy-init]
autorestart=false
autostart=true
command=/etc/init.d/cozy-init
redirect_stderr=true
user=root
priority=3
