# supervisor/cozy-controller.conf
[program:cozy-controller]
autorestart=false
autostart=false
command=cozy-controller
environment=NODE_ENV="production"
redirect_stderr=true
user=root
priority=11
