# supervisor/nginx.conf
[program:nginx]
autorestart=true
autostart=true
command=/usr/sbin/nginx
stdout_events_enabled=true
stderr_events_enabled=true
